# deploy/com.tanya.realestate-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tanya.realestate-bot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/caffeinate</string>
        <string>-i</string>
        <string>-s</string>
        <string>/Users/tanya.vasileva/realestate-bot/.venv/bin/python</string>
        <string>/Users/tanya.vasileva/realestate-bot/bot.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/tanya.vasileva/realestate-bot</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/tanya.vasileva/realestate-bot/bot.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/tanya.vasileva/realestate-bot/bot.log</string>
</dict>
</plist>
